# crontab
@monthly /root/bin/renew-ssl-cert.sh mydomain.com
